# roles/static-page/templates/static-page.service.j2
# {{ ansible_managed }}
# static-page.service — serves the emacboros landing page
# Uses Python's built-in HTTP server (zero dependencies)

[Unit]
Description=emacboros static landing page
After=network.target

[Service]
Type=simple
User={{ static_page_user }}
Group={{ static_page_group }}
WorkingDirectory={{ static_page_root }}
ExecStart=/usr/bin/python3 -m http.server {{ static_page_port }} --bind 10.66.0.3
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths={{ static_page_root }}
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target