# roles/node-exporter/templates/node_exporter.service.j2
# {{ ansible_managed }}
# node_exporter.service — Prometheus host metrics exporter
# Listens on WireGuard IP only — scraped by Prometheus over the mesh.

[Unit]
Description=Prometheus Node Exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ node_exporter_user }}
Group={{ node_exporter_user }}
ExecStart={{ node_exporter_bin }} \
    --web.listen-address={{ wg_ip }}:{{ node_exporter_port }} \
    --web.disable-rootfs-info \
    --collector.filesystem.mount-points-exclude=^/(dev|proc|sys|var/lib/docker/.+)($$|/) \
    --collector.filesystem.fs-types-exclude=^(tmpfs|overlay|squashfs)$

Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelModules=true
ProtectKernelTunables=true

[Install]
WantedBy=multi-user.target