# roles/node-exporter/templates/node_exporter.service.j2
# {{ ansible_managed }}
# node_exporter.service — Prometheus host metrics exporter
# Listens on WireGuard IP only — scraped by Prometheus over the mesh.

[Unit]
Description=Prometheus Node Exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ node_exporter_user }}
Group={{ node_exporter_user }}
ExecStart={{ node_exporter_bin }} \
    --web.listen-address={{ wg_ip }}:{{ node_exporter_port }}

Restart=on-failure
RestartSec=5

# Hardening (relaxed for node_exporter — needs to read /proc, /sys, filesystems)
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true
ProtectKernelModules=true
ProtectKernelTunables=true

[Install]
WantedBy=multi-user.target